# yard/skills/00-web/django_demo/conf/supervisor_server.ini
[program:django-server]
command=python3 manage.py runserver 8081
directory=/home/demo/django_demo
environment=PATH="/home/demo/.venv-hub/py3-django-xadmin/bin"
stdout_logfile=/home/demo/django_demo/log/dj-server.log
stderr_logfile=/home/demo/django_demo/log/dj-server-err.log
autostart=true
autorestart=true
startsecs=10
stopwatisecs=60
priority=998
