# yard/skills/00-dss/celery-demo/conf/supervisor_celery_worker2.ini
[program:django-celery-worker2]
command=celery worker -A 'celery_app' -l INFO -n worker-02@mbp
directory=/home/demo/celery-demo
environment=PATH="/home/demo/.venv-hub/py3-celery4-dj2/bin"
stdout_logfile=/home/demo/celery-demo/log/celery.worker-02.log
stderr_logfile=/home/demo/celery-demo/log/celery-err.worker-02.log
autostart=true
autorestart=true
startsecs=10
stopwatisecs=60
priority=998
